# crontab
0 9-18 * * 1-5 /usr/local/bin/python /usr/src/app/main.py >> /var/log/cron.log 2>&11